# deploy/auto-update.service
[Unit]
Description=Nexus Auto-Update Check

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/nexus/deploy/auto-update.sh
