# deploy/litestream.service
[Unit]
Description=Litestream SQLite Replication
Before=nexus.service

[Service]
EnvironmentFile=/opt/nexus/.env
ExecStart=/usr/bin/litestream replicate -config /etc/litestream.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
