# deploy/todoapp.service
[Unit]
Description=ToDo App (Gunicorn)
After=network.target

[Service]
User=todoapp
Group=todoapp
WorkingDirectory=/opt/todoapp
EnvironmentFile=/opt/todoapp/.env
ExecStart=/opt/todoapp/venv/bin/gunicorn todoapp.wsgi:application --bind 127.0.0.1:8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
